# config/chroot_local-includes/usr/lib/systemd/user/tails-report-squashfs-errors.service
[Unit]
Description=Report errors about disk squashfs being corrupted
ConditionUser=1000

[Service]
Type=oneshot
ExecStart=/usr/local/lib/tails-report-squashfs-errors
RemainAfterExit=yes



